# sqo
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $
#
# /etc/init.d/squash_overlays allows efficient compression of
# Gentoo overlays arborescence
#
# It requires support for the loop device and squashfs enabled in the kernel,
# module autoloading is also *highly* recommended.
# sys-fs/squashfs and sys-fs/unionfs are necessary for read-write support.
#
# Author: Mathias Laurin <mathias_laurin@users.sourceforge.net>
# 2006-11-28, v.0.1.5(4)
# patched for aufs - yoshi314

PORTDIR=/usr/overlays
SQFS_DIRNAME=/mnt/storage
SQFS_CUR="$SQFS_DIRNAME/overlays.sqfs"
SQFS_NEW="$SQFS_DIRNAME/overlays-current.sqfs"
SQFS_OLD="$SQFS_DIRNAME/overlays-old.sqfs"
DEF_RW="/dev/shm/.overlays-rw"

PMANAGER=paludis
FS=unionfs

#compression: currently available are xz, lzo, gzip
COMP=gzip

depend() {
	use local
}

if [ "$PMANAGER" = "paludis" ] ; then
	_UID=$(grep ^paludisbuild /etc/passwd | cut -d\: -f3)
	_GID=$(grep ^paludisbuild /etc/group | cut -d\: -f3)
else
	_UID=$(grep ^overlays /etc/passwd | cut -d\: -f3)
	_GID=$(grep ^overlays /etc/group | cut -d\: -f3)
fi


setup() {
	rm "$SQFS_OLD" "$SQFS_NEW"
	einfo "Attempting to build initial squashfs image from $PORTDIR"
	mksquashfs $PORTDIR $SQFS_NEW -force-uid ${_UID} -force-gid ${_GID} -comp $COMP -no-recovery 2>/dev/null
}

start() {
   ebegin "Mounting read-only squashfs image"

   if [ ! -f "$SQFS_CUR" ] ; then
	if [ ! -f "$SQFS_OLD" ]; then
		einfo "it appears you have no squashfs images yet. If you want to create initial images from $PORTDIR, run $0 setup"
		retval=1
		eend $retval
	fi
   fi

   mount -rt squashfs -o loop,nodev,noexec $SQFS_CUR $PORTDIR
   retval=$?
   eend $retval
   [ $retval -ne 0 ] && return $retval

   ebegin "Mounting read-write with unionfs"
   if [ ! $PORTAGE_RW ]
   then
      einfo "  mounted in tmpfs (RAM)"
      PORTAGE_RW="${DEF_RW}"
   fi
   [ -d $PORTAGE_RW ] || mkdir -p $PORTAGE_RW
   chmod 0750 $PORTAGE_RW
   chown ${_UID}:${_GID} $PORTAGE_RW

   if [ "$FS" = "aufs" ] ; then
     mount -t aufs -o br=$PORTAGE_RW=rw:$PORTDIR=ro none $PORTDIR
   else
     mount -t unionfs -o nodev,noexec,dirs=$PORTAGE_RW=rw:$PORTDIR=ro unionfs $PORTDIR 
   fi

   eend $?
}

stop() {
   ebegin "Updating overlays tree"
   [ ! $PORTAGE_RW ] && PORTAGE_RW="${DEF_RW}"
   if [ ! -z `ls -A --ignore=.wh* $PORTAGE_RW | head -n1` ]
   then
      einfo "  Syncing the tree"
      mv -f $SQFS_NEW $SQFS_OLD
      mksquashfs $PORTDIR $SQFS_NEW -force-uid ${_UID} -force-gid ${_GID} -comp $COMP -no-recovery 2>/dev/null
      retval=$?
      ln -sf $SQFS_NEW $SQFS_CUR
   else
      einfo "  Nothing to do"
      retval=0
   fi
   eend $retval

   ebegin "Unmounting the tree"
   umount -t ${FS} $PORTDIR
   umount -t squashfs $PORTDIR
   rm -rf $PORTAGE_RW
   eend 0
} 
